# .config/systemd/user/tmux@scratch.service.d/override-dir.conf
[Service]
ExecStart=
ExecStart=/usr/bin/tmux -f ${HOME}/.config/tmux/config new-session -ds scratch -c ${HOME}/scratch
